# guitar-arpeggiator.service
[Unit]
Description=Guitar Arpeggiator Headless Mode
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=audio
WorkingDirectory=/home/pi/guitar-arpeggiator
Environment=PATH=/home/pi/guitar-arpeggiator/venv/bin
ExecStart=/home/pi/guitar-arpeggiator/venv/bin/python headless_mode.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# GPIO access
SupplementaryGroups=gpio

# Audio access
SupplementaryGroups=audio

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
